# move2docker/etc/service/birdlg_proxy/run
#!/bin/bash
. /.denv
if [[ "$BIRDLG_PROXY" == 1 ]]; then
  echo "birdlg_proxy enabled, start"
  exec bird-lg-proxy --listen=:$BIRDLG_PROXY_PORT
else
  echo "birdlg_proxy not enabled, down"
  sv down birdlg_proxy
  exec sleep infinity
fi
