# move2docker/etc/service/crond/run
#!/bin/bash
echo "crond up"
exec 2>&1
cd /
pkill cron
exec cron -f
